# guacamole/rootfs/etc/services.d/tomcat/run
#!/usr/bin/with-contenv bashio
export GUACAMOLE_HOME=/data/guacamole
# Ensure Tomcat temp dirs exist
mkdir -p /run/tomcat
exec /opt/tomcat/bin/catalina.sh run
