# guacamole/rootfs/etc/services.d/guacd/run
#!/usr/bin/with-contenv bashio
exec guacd -f
