# guacamole/rootfs/etc/services.d/guacd/run
#!/usr/bin/with-contenv sh
# guacd is provided by guacamole-server package
# Bind explicitly to IPv4 localhost to match guacamole.properties
exec /usr/sbin/guacd -f -b 127.0.0.1 -l 4822
