# guacamole/rootfs/etc/services.d/guacd/run
#!/usr/bin/with-contenv bashio
# guacd is provided by guacamole-server package
exec /usr/sbin/guacd -f
