# guacamole/rootfs/etc/services.d/guacd/run
#!/usr/bin/with-contenv sh
# Start guacd, locating the binary regardless of install prefix
set -e
GUACD_BIN=$(command -v guacd || true)
if [ -z "$GUACD_BIN" ]; then
	if [ -x /usr/sbin/guacd ]; then GUACD_BIN=/usr/sbin/guacd; fi
fi
if [ -z "$GUACD_BIN" ]; then
	if [ -x /usr/local/sbin/guacd ]; then GUACD_BIN=/usr/local/sbin/guacd; fi
fi
if [ -z "$GUACD_BIN" ]; then
	echo "[guacd] ERROR: guacd binary not found" >&2
	exit 1
fi

# Start guacd with info log level; rely on default plugin path
exec "$GUACD_BIN" -f -b 127.0.0.1 -l 4822 -L info
